# etc/system.d/system/monitoring_trigger.path
[Unit]
Description=Watch for first user login.

[Path]
PathExistsGlob=/run/systemd/users/*

[Install]
WantedBy=multi-user.target
