# etc/profile.d/show-monitoring-once.sh
#!/bin/bash
/usr/local/libexec/monitoring.sh
