# etc/systemd/system/monitoring_broadcast.service
[Unit]
Description=Runs monitoring_broadcast.sh

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/monitoring_broadcast.sh
